# server/deploy/wanderclaw.service
[Unit]
Description=WanderClaw Onboarding + Chat Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace-pm/projects/wanderClaw/onboarding
ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=5
Environment=WANDERCLAW_USERS_BASE=/root/wanderClaw/users
Environment=OPENCLAW_CONFIG=/root/.openclaw/openclaw.json
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wanderclaw

[Install]
WantedBy=multi-user.target
